# service/readImageSensor.timer
[Unit]
Description=timer read sensor and shoot image readImageSensor.py

[Timer]
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
